# systemd/user/fusuma.service
[Unit]
Description=Fusuma - gesture detection

[Service]
ExecStart=/usr/bin/fusuma &> 
LogLevelMax=3

[Install]
WantedBy=default.target
